# install/systemd.service.tmpl
[Unit]
Description=Orbit autonomous coding loop (__TARGET__)
After=network-online.target

[Service]
Type=simple
Environment=ORBIT_HOME=__HOME__
Environment=ORBIT_TARGET=__TARGET__
ExecStart=/bin/bash __RUN__ __TARGET__
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
